# poky/meta-myproject/recipes-hello/source/hello.service
[Unit]
Description=Hello World startup script

[Service]
ExecStart=/usr/bin/hello

[Install]
WantedBy=multi-user.target

